# bash_lab/gunicorn.service
[Unit]
Description=Run the api and flask app

[Service]
ExecStart=/path/to/main-service

[Install]
WantedBy=multi-user.target
